# ansible/roles/liberty/templates/liberty.service.j2
[Unit]
Description=Open Liberty Server - {{ liberty_server_name }}
Documentation=https://openliberty.io/docs/
After=network.target

[Service]
Type=forking
User={{ liberty_user }}
Group={{ liberty_group }}
Environment="JAVA_HOME=/usr/lib/jvm/java-17-amazon-corretto"
Environment="WLP_USER_DIR={{ liberty_wlp_dir }}/usr"
Environment="LOG_DIR=/var/log/liberty"

ExecStart={{ liberty_wlp_dir }}/bin/server start {{ liberty_server_name }}
ExecStop={{ liberty_wlp_dir }}/bin/server stop {{ liberty_server_name }}
ExecReload={{ liberty_wlp_dir }}/bin/server stop {{ liberty_server_name }} && {{ liberty_wlp_dir }}/bin/server start {{ liberty_server_name }}

PIDFile={{ liberty_servers_dir }}/{{ liberty_server_name }}/workarea/.pid/{{ liberty_server_name }}.pid

Restart=on-failure
RestartSec=10

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ liberty_install_dir }} /var/log/liberty

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

[Install]
WantedBy=multi-user.target
